# phpmyadmin/rootfs/etc/s6-overlay/s6-rc.d/init-phpmyadmin/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: phpMyAdmin
# Configures phpMyAdmin
# ==============================================================================
declare host
declare password
declare port
declare username
declare database

if ! bashio::fs.file_exists "/data/config.secret.inc.php"; then
    cat > /data/config.secret.inc.php <<EOT
<?php
\$cfg['blowfish_secret'] = '$(tr -dc 'a-zA-Z0-9~!@#$%^&*_()+}{?></";.,[]=-' < /dev/urandom | fold -w 32 | head -n 1)';
EOT
fi
chmod 644 /data/config.secret.inc.php

if bashio::config.has_value 'host' ; then
  host=$(bashio::config 'host')
  if bashio::config.has_value 'port' ; then
    port=$(bashio::config 'port')
  else
    port=3306
  fi
  username=$(bashio::config 'username')
  password=$(bashio::config 'password')
else
  host=$(bashio::services "mysql" "host")
  port=$(bashio::services "mysql" "port")
  username=$(bashio::services "mysql" "username")
  password=$(bashio::services "mysql" "password")
fi

database=$(\
    mariadb \
        -u "${username}" -p"${password}" \
        -h "${host}" -P "${port}" \
        --skip-column-names \
        -e "SHOW DATABASES LIKE 'phpmyadmin';"
)

if ! bashio::var.has_value "${database}"; then
    bashio::log.info "Creating database for phpMyAdmin"
    mariadb \
        -u "${username}" -p"${password}" \
        -h "${host}" -P "${port}" \
            < /var/www/phpmyadmin/sql/create_tables.sql
fi
